# deploy/huifu-server.service
[Unit]
Description=惠福星链 · 全病程健康协同平台
After=network.target mysql.service redis.service

[Service]
Type=simple
User=huifu
Group=huifu
WorkingDirectory=/opt/huifu

EnvironmentFile=/opt/huifu/.env
ExecStart=/opt/huifu/huifu-server
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/huifu /data/huifu /var/log/huifu

# Resource limits
MemoryLimit=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=huifu-server

[Install]
WantedBy=multi-user.target
